# project/backend/blockchain/tasks.crontab
# Updating index and images with 1 minute delay between
20 4 * * * /blockchain/blago &>> /var/log/blago.log
